# deploy/matibot-v7-promotion.service
[Unit]
Description=MatiTradingBot v7 promotion controller (fail closed)
After=matibot-v7-shadow.service

[Service]
Type=simple
User=__RUN_USER__
WorkingDirectory=__APP_DIR__
Environment=PYTHONUNBUFFERED=1
ExecStart=__APP_DIR__/.venv/bin/python tools/v7_promotion_controller.py --watch --promote
Restart=always
RestartSec=30
NoNewPrivileges=true
PrivateTmp=true
MemoryMax=150M

[Install]
WantedBy=multi-user.target
